# infra/ansible/roles/observability/templates/node_exporter.service.j2
# Copyright © 2026 mindicator & silicon bags quartet.
# SPDX-License-Identifier: AGPL-3.0-or-later
# This file is part of Mycelium, licensed under the GNU Affero General Public License v3.0 or
# later. See the LICENSE file in the repository root.
#
# Mycelium Phase 0 — hardened systemd unit for node_exporter (rendered by the observability role).
# Author: mindicator & silicon bags quartet
#
# Bound to {{ node_exporter_listen }} (loopback). Not exposed by the host firewall.
[Unit]
Description=Mycelium node_exporter (host metrics, loopback only)
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_group }}
ExecStart={{ node_exporter_bin_path }} --web.listen-address={{ node_exporter_listen }}
Restart=on-failure
RestartSec=5s

# --- Privilege containment -------------------------------------------------
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=

# --- Filesystem / kernel isolation ----------------------------------------
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
PrivateDevices=true
DevicePolicy=closed
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
UMask=0077

# --- Network / syscall surface --------------------------------------------
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete

[Install]
WantedBy=multi-user.target
